# silex3.service
[Unit]
Description=Silex3 Actuator Control Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=axton
WorkingDirectory=/home/axton/Silex3
Environment="PATH=/home/axton/Silex3/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/axton/Silex3/venv/bin/python3 /home/axton/Silex3/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

